# supervisord.conf
[supervisord]
nodaemon=true

[program:fluentd]
command=/usr/local/bin/fluentd -c /etc/fluent/fluent.conf --use-v1-config

[program:ngnix]
command=nginx

[program:cron]
command=crond -f